# docker/cron/crontab
* * * * * sh php bin/console sylius:remove-expired-carts
* * * * * sh php bin/console sylius:sylius:cancel-unpaid-orders
